# platform/marvell/sonic-platform-db98cx8514-10cc/debian/sonic-platform-db98cx8514-10cc-db98cx8514.postinst
#!/bin/sh
# postinst script for sonic-platform-db98cx8540-16cd
#
# see: dh_installdeb(1)

chmod a+x /usr/local/bin/db98cx8514-10cc-init.sh
systemctl enable db98cx8514-10cc-init.service
systemctl start db98cx8514-10cc-init.service
exit 0

