# platform/marvell-arm64/sonic-platform-tg4810m/debian/sonic-platform-tg4810m.postinst
#!/bin/sh
# postinst script for sonic-platform-tg4810m
#
# see: dh_installdeb(1)

sh /usr/sbin/tg4810m_plt_setup.sh
chmod a+x /usr/local/bin/tg4810minit.sh
systemctl enable tg4810minit.service
systemctl start tg4810minit.service
exit 0

